# conf/etc/supervisor/conf.d/supervisord_npm.conf
[program:npm]
directory=/opt/framework
command=sh -c 'NODE_ENV=dev npm run watch'
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
